# docker-ray/pd2jira/apache-php/sv/log/run
#!/bin/sh
service="pd2jira"
logdir="/var/log/${service}"
mkdir -p ${logdir}
exec 2>&1
exec /usr/bin/svlogd ${logdir}